# Chapter10/06-Sending Docker logs to AWS CloudWatch Logs/etc/systemd/system/docker.service.d/aws.conf
[Service]
Environment="AWS_ACCESS_KEY_ID=AKIAJ..."
Environment="AWS_SECRET_ACCESS_KEY=SW+jdHKd.."
